# Docker/Kubernetes/KubernetesCluster/package/Resources/scripts/systemd/confd.service
[Unit]
Description=Configuration Service
Requires=etcd.service

[Service]
Restart=always
ExecStart=/usr/local/bin/confd \
 -backend etcd \
 -node %%MASTER_NODE%%:4001 \
 -interval 15

[Install]
WantedBy=basic.target
